# ldap-oidc-sync.timer
[Unit]
Description=Sync LDAP/OIDC Users for Proxmox

[Timer]
OnCalendar=*-*-* 4:00:00
Persistent=true

[Install]
WantedBy=timers.target